# scripts/nuc/templates/quiz-nuc-restart.timer
[Unit]
Description=shh — redémarrage nocturne Chromium (purge freeze potentiels hors séances)

[Timer]
OnCalendar=*-*-* 04:00:00
Persistent=true
Unit=quiz-nuc-restart.service

[Install]
WantedBy=timers.target
