# scripts/nuc/templates/quiz-nuc-restart.service
[Unit]
Description=shh — redémarrage préventif du kiosque Chromium (oneshot)

[Service]
Type=oneshot
ExecStart=/usr/bin/systemctl restart quiz-nuc-chromium.service
